# BeAST_Pi_Code/beast-synthetic.service
[Unit]
Description=BeAST Synthetic Data Playback
After=network.target beast-voice.service
Wants=beast-voice.service

[Service]
Type=simple
EnvironmentFile=/etc/beast/environment
User=${BEAST_USER}
WorkingDirectory=${BEAST_HOME}/Beast/BeAST_Pi_Code/Live Connections Simulator Test
Environment="PATH=${BEAST_HOME}/.beast-venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=${BEAST_HOME}/.beast-venv/bin/python beast_synthetic_playback.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
