# BeAST_Pi_Code/beast-voice.service
[Unit]
Description=BeAST Voice Assistant
After=network.target sound.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
EnvironmentFile=/etc/beast/environment
User=${BEAST_USER}
WorkingDirectory=${BEAST_HOME}/Beast/BeAST_Pi_Code/TheBeast/TheBeast/beast
Environment="PATH=${BEAST_HOME}/.beast-venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="HF_HUB_OFFLINE=1"
Environment="TRANSFORMERS_OFFLINE=1"
Environment="BEAST_DB_PASSWORD=beast"
ExecStart=${BEAST_HOME}/.beast-venv/bin/python main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
