# BeAST_Pi_Code/beast-synthetic.service
[Unit]
Description=BeAST Synthetic Data Playback
After=network.target beast-voice@%i.service
Wants=beast-voice@%i.service

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/Beast/BeAST_Pi_Code/Live Connections Simulator Test
Environment="PATH=/home/%i/.beast-venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/%i/.beast-venv/bin/python beast_synthetic_playback.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
