# BeAST_Pi_Code/beast-synthetic.service
[Unit]
Description=BeAST Synthetic Data Playback
After=network.target beast-voice.service
Wants=beast-voice.service

[Service]
Type=simple
User=%u
WorkingDirectory=%h/Beast/BeAST_Pi_Code/Live Connections Simulator Test
Environment="PATH=%h/.beast-venv/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=%h/.beast-venv/bin/python beast_synthetic_playback.py
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
