# BeAST_Pi_Code/TheBeast/TheBeast/beast/beast-autostart.service
[Unit]
Description=BeAST Wearable AI System
After=network.target sound.target
Wants=sound.target

[Service]
Type=simple
User=beast1
Group=beast1
WorkingDirectory=/home/beast1/Beast/BeAST_Pi_Code/TheBeast/TheBeast/beast

Environment="PATH=/home/beast1/beast-venv/bin:/usr/local/bin:/usr/bin:/bin"
Environment="BEAST_HOME=/home/beast1/Beast/BeAST_Pi_Code/TheBeast/TheBeast/beast"
Environment="VIRTUAL_ENV=/home/beast1/beast-venv"

ExecStartPre=/bin/sleep 5
ExecStart=/home/beast1/beast-venv/bin/python /home/beast1/Beast/BeAST_Pi_Code/TheBeast/TheBeast/beast/main.py

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
